# deploy/systemd/freya.service
[Unit]
Description=Freya API Server
After=network.target

[Service]
Type=simple
User=freya
WorkingDirectory=/opt/freya
ExecStart=/opt/freya/.venv/bin/freya serve --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
Environment=HOME=/opt/freya
# Binding 0.0.0.0 requires authentication. This file MUST exist and contain:
#   FREYA_API_KEY=<key>      (generate one: `freya auth generate-key`)
# It is not prefixed with "-", so the unit fails to start if the file is
# missing — preventing an accidentally unauthenticated public server.
EnvironmentFile=/etc/freya/env

[Install]
WantedBy=multi-user.target
